# ocr.service
[Unit]
Description=tesseract ocr

[Service]
User=scan
ExecStart=/usr/local/bin/ocr.sh

[Install]
WantedBy=multi-user.target
